# package/solaris11/ercole-agent.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle
  SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    Manifest created by svcbundle (2018-ott-04 14:33:54+0200)
-->
<service_bundle type="manifest" name="ercole-agent">
    <service version="1" type="service" name="ercole-agent">
        <dependency restart_on="none" type="service"
            name="multi_user_dependency" grouping="require_all">
            <service_fmri value="svc:/milestone/multi-user"/>
        </dependency>
        <exec_method timeout_seconds="180" type="method" name="start"
            exec="/lib/svc/method/ercole-agent-start"/>
        <exec_method timeout_seconds="60" type="method" name="stop"
            exec=":kill"/>
        <exec_method timeout_seconds="60" type="method" name="refresh"
            exec=":true"/>
        <instance enabled="true" name="default"/>
        <template>
            <common_name>
                <loctext xml:lang="C">
                        ercole-agent
                </loctext>
            </common_name>
            <description>
                <loctext xml:lang="C">
                        The ercole-agent service.
                </loctext>
            </description>
        </template>
    </service>
</service_bundle>
